# systemd/codehero-daemon.service
[Unit]
Description=CodeHero Daemon
After=network.target mysql.service

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/opt/codehero
ExecStart=/usr/bin/python3 /opt/codehero/scripts/claude-daemon.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/codehero/daemon.log
StandardError=append:/var/log/codehero/daemon.log

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
